# nfsroot1/etc/inittab
ttySTM0::askfirst:/home/ahmed/embedded-linux-labs/tinysystem/nfsroot/bin/sh

